# repo/sysklogd/sysklogd.initd
#!/sbin/openrc-run

command=/bin/syslogd
